# recipes/files/etc/systemd/system/apache2.service.d/override.conf
# Override the default PrivateTmp=true
[Service]
PrivateTmp=false
