# recipes/files/etc/systemd/system/apache2.service.d/override.conf
[Service]
PrivateTmp=false
Restart=always
