# recipes/files/etc/systemd/system/git-sync-check.timer
[Unit]
Description=Timer for git-sync-check.service

[Timer]
OnCalendar=*:0/10
Unit=git-sync-check.service

[Install]
WantedBy=timers.target
