# recipes/files/etc/systemd/system/datadog-agent.service.d/override.conf
[Service]
# Enable Datadog log collection without editing the secret-bearing datadog.yaml.
# Equivalent to logs_enabled: true.
Environment=DD_LOGS_ENABLED=true
